# systemd/statusframe.service
[Unit]
Description=StatusFrame Application Monitoring Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/statusframe
Environment="PATH=/usr/local/go/bin:/usr/bin:/bin"
EnvironmentFile=/home/ubuntu/statusframe/.env
ExecStart=/home/ubuntu/statusframe/statusframe
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=statusframe

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/statusframe

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
